# robot_voice/systemd/ollama-robot.service
[Unit]
Description=Ollama service for Raspberry Pi robot command parsing
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=phuong
Group=phuong
Environment=OLLAMA_HOST=127.0.0.1:11434
Environment=OLLAMA_KEEP_ALIVE=30m
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
